# tools/wsl_ubuntu/system/etc/profile.d/nvm_autoload.sh
#! /bin/sh

# Enable nvm if available
if [ -f ~/.nvm/nvm.sh ]; then
  . ~/.nvm/nvm.sh
fi
