# tools/wsl_ubuntu/system/etc/profile.d/pyenv_autoload.sh
#! /bin/sh

# Enable pyenv if available
if [ -d "$HOME/.pyenv" ]; then
  export PYENV_ROOT="$HOME/.pyenv"
  command -v pyenv >/dev/null || export PATH="$PYENV_ROOT/bin:$PATH"
  eval "$(pyenv init -)"
  eval "$(pyenv virtualenv-init -)"
fi
